# cron.d/videk-hosts
*/5 * * * *     /usr/bin/python /root/videk-hosts/videk-hosts.py
