# cron.d/videk-hosts
*/5 * * * *     /usr/bin/python /home/videk-hosts/videk-hosts.py
